# overlay/etc/sv/nginx/run
#!/bin/bash

export TZ="${NAGIOS_TIMEZONE}"

exec nginx -g daemon off;
